# hardware/pi/onelink-brain.service
[Unit]
Description=OneLink Brain (MQTT to WebSocket bridge for kiosk)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Onelink/hardware/pi
# Prefer /etc/onelink/brain.env — set MQTT_BROKER=broker.emqx.io there.
EnvironmentFile=-/etc/onelink/brain.env
Environment=MQTT_BROKER=broker.emqx.io
Environment=MQTT_PORT=1883
Environment=MQTT_TOPIC=onelink/hardware/tap
Environment=MQTT_TOPIC_HEARTBEAT=onelink/hardware/heartbeat
Environment=KIOSK_WS_HOST=0.0.0.0
Environment=KIOSK_WS_PORT=8765
Environment=DEFAULT_TERMINAL=main_kiosk
Environment=AUTO_PAY=0
# Use the Pi venv that already has paho-mqtt + websockets.
ExecStart=/home/pi/onelink-venv/bin/python /home/pi/Onelink/hardware/pi/brain.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
